# files/pbs-backup.timer
[Unit]
Description=Daily PBS backup

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
